# cde/admin/BuildTools/tog/cde.crontab
###############################################################################
#
# Builds and reports 
#
15 * * * *  /project/dt/scripts/cron_scripts -db /project/dt/admin/cron/cde.crondb
